# sys-apps/world/files/world.timer
[Unit]
Description=Daily Cleanup of User's Temporary Directories

[Timer]
OnStartupSec=1min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
